# Library/LaunchAgents/com.1password.SSH_AUTH_SOCK.plist
/* Per https://developer.1password.com/docs/ssh/agent/compatibility/ */
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.1password.SSH_AUTH_SOCK</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/sh</string>
    <string>-c</string>
    <string>for s in "$HOME"/Library/Group\ Containers/*.com.1password/t/agent.sock; do [ -S "$s" ] &amp;&amp; /bin/ln -sf "$s" "$HOME/.ssh/1password.agent.sock" &amp;&amp; break; done</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
</dict>
</plist>
